# deploy/market-listen.service
# systemd unit for the --listen long-polling daemon (near-instant on-demand digests).
# Install:
#   sudo cp deploy/market-listen.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now market-listen
# Notes:
# - The env file uses `export KEY=val` lines, which systemd's EnvironmentFile cannot parse,
#   so we source it in a shell and `exec` python3 — exec replaces the shell, letting systemd
#   supervise (signal/restart) the real Python process rather than a bash wrapper.
# - Replace /home/youruser with your actual home directory.

[Unit]
Description=market_alerts Telegram long-poll listener
After=network-online.target
Wants=network-online.target

[Service]
User=youruser
WorkingDirectory=/home/youruser/market-extremes-alerter
ExecStart=/bin/bash -c '. /home/youruser/.env && exec /usr/bin/python3 /home/youruser/market-extremes-alerter/market_alerts.py --listen'
Restart=always
RestartSec=5
StandardOutput=append:/home/youruser/market-extremes-alerter/run.log
StandardError=append:/home/youruser/market-extremes-alerter/run.log

[Install]
WantedBy=multi-user.target
